# macro/littrack/nightly_tests/tests_cdash.cron
0 1 * * * `. /data.local1/andrey/cdash_tests/cbmroot/macro/littrack/nightly_tests/run_test_cdash.sh`
