# macro/littrack/nightly_tests/tests.cron
0 23 * * * `. /u/andrey/cbm/trunk/cbmroot/macro/littrack/nightly_tests/run_test.sh cdash`
